# roles/lisk-node/templates/profile.d/lisk_network.sh
export LISK_NETWORK={{ lisk_node_network }}
